# files/etc/supervisor/conf.d/postgresql.conf
[program:postgresql]
command=run_postgres
stdout_events_enabled=true
stderr_events_enabled=true
